# supervisord-redis.conf
[program:redis]
command=/start-redis.sh
numprocs=1
autostart=true
autorestart=true