# supervisord-mongodb.conf
[program:mongodb]
command=/start-mongodb.sh
numprocs=1
autostart=true
autorestart=true